# support/supervisord.ini
[inet_http_server]
port = :9001

[program:micbot]
directory = /root
command = ./micbot_arm64
stdout_logfile = /root/logs/stdout.log
stdout_logfile_maxbytes = 50MB
stdout_logfile_backups = 2
stderr_logfile = /root/logs/stderr.log
stderr_logfile_maxbytes = 50MB
stderr_logfile_backups = 2

